# data/passwordchecker.service
[Unit]
Description=Checking password expiration

[Service]
Type=simple
ExecStart=/usr/bin/passwordchecker

[Install]
WantedBy=multi-user.target